# system-update.cron
#!/bin/bash
# Update the server
0 22 * 0 * root yum update -y
